# systemd/user/gammastep-stop.service
[Unit]
Description=Stop gammastep colour adjustment

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user stop gammastep.service
